# systemd/hsd.service
[Unit]
Description=hsd mainnet service

[Service]
Type

[Install]
WantedBy=multi-user.target
